# aws/packer/setup/e2b-api.service
[Unit]
Description=E2B API Server — REST API for sandbox management
After=network-online.target docker.service e2b-orchestrator.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/opt/e2b/api.env
ExecStart=/usr/local/bin/e2b-api
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/api.log
StandardError=append:/var/log/api.log

# Hardening
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
